# rootfs/etc/s6/postgres/run
#!/bin/sh

set -eux

su-exec app postgres
