# rootfs/etc/s6/gogs/finish
#!/bin/sh

set -eux

s6-svscanctl -t /etc/s6
